# udev-video-post-example.rules
ACTION=="add", SUBSYSTEM=="pci", ATTRS{class}=="0x030000", RUN+="/usr/sbin/vbetool udevpost %k"
